# cron/crontab
0 * * * * /cron/cronjob.sh >> /var/log/cronjob.log 2>&1
