# cron/crontab
0 * * * * /cron/job.sh >> /var/log/cronjob.log 2>&1
